# supervisord.conf
; CONFIGURACION INICIAL problema con /tmp
; [unix_http_server]
; file=/tmp/supervisor.sock 

; [supervisord]
; logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
; pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
; user=leninelio

; [rpcinterface:supervisor]
; supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; [supervisorctl]
; serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

; [program:gunicorn]
; command=/home/leninelio/django_test_gcp/venv/bin/gunicorn -w 3 backend.wsgi
; directory=/home/leninelio/django_test_gcp
; user=leninelio
; autostart=true
; autorestart=true
; redirect_stderr=true
; environment=DJANGO_SETTINGS_MODULE="backend.settings"



; Permisos necesarios
; sudo chown -R leninelio:leninelio /home/leninelio/django_test_gcp
; sudo chmod 755 /home/leninelio/django_test_gcp

[unix_http_server]
file=/home/leninelio/django_test_gcp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/home/leninelio/django_test_gcp/supervisord.log  ; (main log file; default $CWD/supervisord.log)
pidfile=/home/leninelio/django_test_gcp/supervisord.pid  ; (supervisord pidfile; default supervisord.pid)
user=leninelio  ; (add this line)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/leninelio/django_test_gcp/supervisor.sock  ; (use a unix:// URL  for a unix socket)

[program:gunicorn]
command=/home/leninelio/django_test_gcp/venv/bin/gunicorn -w 3 backend.wsgi
directory=/home/leninelio/django_test_gcp
user=leninelio
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/leninelio/django_test_gcp/gunicorn_stdout.log
stderr_logfile=/home/leninelio/django_test_gcp/gunicorn_stderr.log
environment=DJANGO_SETTINGS_MODULE="backend.settings"
